# rootfs/etc/services.d/jackett/run
#!/usr/bin/execlineb -P

with-contenv
importas PUID PUID
importas PGID PGID
posix-cd /Jackett
s6-setuidgid ${PUID}:${PGID}
./jackett --NoUpdates
